# src/system/etc/systemd/system/docker.socket
[Install]
WantedBy=sockets.target

[Socket]
ListenStream=/var/run/docker.sock
SocketGroup=docker
SocketMode=0660
SocketUser=root

[Unit]
Description=docker socket
PartOf=docker.service
