# src/system/etc/systemd/system/powertop.service
[Install]
WantedBy=multi-user.target

[Service]
RemainAfterExit=true
Type=oneshot

ExecStart=/usr/bin/powertop \
    --auto-tune
