# src/system/etc/NetworkManager/dispatcher.d/99rawdns
#!/bin/bash

IF="$1"
STATUS="$2"

if [[ ! "${IF}" =~ veth* ]] && ( [[ "${STATUS}" == "down" ]] || [[ "${STATUS}" == "up" ]] ) ; then
    systemctl restart rawdns
fi
