# src/system/etc/rc.local
#!/bin/bash

powertop --auto-tune

exit 0
