# src/system/etc/rc.local
#!/bin/bash
set -e

powertop --auto-tune

exit 0
