# src/system/etc/rc.local
#!/bin/sh

powertop --auto-tune

exit 0
